# ansible/roles/os-auto-update/templates/os-update-notify.service.j2
[Unit]
Description=OS auto-update result notification

[Service]
Type=oneshot
ExecStart=/usr/local/bin/os-update-notify.sh
