# raspberry-pi/recovery/recovery.service
[Unit]
Description=K3s Recovery Webhook Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Infisical から全シークレットを注入してサービスを起動する
# /opt/recovery/.infisical-auth には Infisical 認証情報のみ置く (シークレット本体は置かない)
#   INFISICAL_CLIENT_ID=...
#   INFISICAL_CLIENT_SECRET=...
#   INFISICAL_PROJECT_ID=...
# chmod 600 /opt/recovery/.infisical-auth で保護すること
EnvironmentFile=/opt/recovery/.infisical-auth
ExecStart=infisical run -- /usr/bin/python3 /opt/recovery/recover.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
